# home/config/autostart/auto-download-dir.desktop
[Desktop Entry]
Exec=/usr/bin/mkdir -p /dev/shm/下载
Icon=
Name=自动创建下载文件夹
Path=
Terminal=False
Type=Application
