# home/config/autostart/latte-dock-wrapper.desktop
[Desktop Entry]
Exec=/home/yidaozhan/.local/bin/latte-dock-wrapper
Icon=
Name=latte-dock-wrapper
Path=
Terminal=False
Type=Application
